# security-system.service
# systemd unit template — do NOT edit this file directly.
# Run setup_service.sh to generate and install the real unit file.
[Unit]
Description=Smart Security System (Flask)
After=network.target mysql.service mariadb.service
Wants=network.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__PARENT_DIR__
Environment=PYTHONPATH=__PARENT_DIR__
ExecStart=__VENV_PYTHON__ -m __PACKAGE_NAME__.run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
